# jboss-7.1-standalone/roles/jboss-standalone/files/jboss-standalone.service
[Service]
ExecStart=/etc/rc.d/init.d/jboss start
Restart=on-failure
